# deploy/systemd/worldcup.service
[Unit]
Description=World Cup Boys API and SPA
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/World-Cup-Project
EnvironmentFile=/home/ubuntu/World-Cup-Project/.env
# Use system node + tsx cli.mjs (reliable under systemd; do not use npm run server).
ExecStart=/usr/bin/node /home/ubuntu/World-Cup-Project/node_modules/tsx/dist/cli.mjs src/server/index.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
