# deploy/systemd/worldcup.service
[Unit]
Description=World Cup Boys API and SPA
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/World-Cup-Project
EnvironmentFile=/home/ubuntu/World-Cup-Project/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
