# deploy/systemd/worldcup-monitor.timer
[Unit]
Description=Check live site health every 2 minutes and auto-recover on outage

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
