# deploy/systemd/worldcup.service
[Unit]
Description=World Cup Boys API and SPA
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/World-Cup-Project
EnvironmentFile=/home/ubuntu/World-Cup-Project/.env
# NODE_ENV and secrets come from .env only — do not set NODE_ENV here (npm/tsx path issues).
ExecStart=/home/ubuntu/World-Cup-Project/node_modules/.bin/tsx src/server/index.ts
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
