# deploy/systemd/worldcup-jobs.service
[Unit]
Description=World Cup Boys background jobs (locks + football-data sync)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/World-Cup-Project
EnvironmentFile=/home/ubuntu/World-Cup-Project/.env
ExecStart=/home/ubuntu/World-Cup-Project/node_modules/.bin/tsx src/server/jobs.ts
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
